# var/lib/pacman/local/colord-1.2.9-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.freedesktop.ColorManager.conf
usr/
usr/bin/
usr/bin/cd-create-profile
usr/bin/cd-fix-profile
usr/bin/cd-iccdump
usr/bin/cd-it8
usr/bin/colormgr
usr/include/
usr/include/colord-1/
usr/include/colord-1/colord-private.h
usr/include/colord-1/colord-session/
usr/include/colord-1/colord-session/cd-session.h
usr/include/colord-1/colord.h
usr/include/colord-1/colord/
usr/include/colord-1/colord/cd-buffer.h
usr/include/colord-1/colord/cd-client-sync.h
usr/include/colord-1/colord/cd-client.h
usr/include/colord-1/colord/cd-color.h
usr/include/colord-1/colord/cd-device-sync.h
usr/include/colord-1/colord/cd-device.h
usr/include/colord-1/colord/cd-dom.h
usr/include/colord-1/colord/cd-edid.h
usr/include/colord-1/colord/cd-enum.h
usr/include/colord-1/colord/cd-icc-store.h
usr/include/colord-1/colord/cd-icc-utils.h
usr/include/colord-1/colord/cd-icc.h
usr/include/colord-1/colord/cd-interp-akima.h
usr/include/colord-1/colord/cd-interp-linear.h
usr/include/colord-1/colord/cd-interp.h
usr/include/colord-1/colord/cd-it8-utils.h
usr/include/colord-1/colord/cd-it8.h
usr/include/colord-1/colord/cd-math.h
usr/include/colord-1/colord/cd-profile-sync.h
usr/include/colord-1/colord/cd-profile.h
usr/include/colord-1/colord/cd-quirk.h
usr/include/colord-1/colord/cd-sensor-sync.h
usr/include/colord-1/colord/cd-sensor.h
usr/include/colord-1/colord/cd-spectrum.h
usr/include/colord-1/colord/cd-transform.h
usr/include/colord-1/colord/cd-version.h
usr/include/colord-1/colorhug/
usr/include/colord-1/colorhug/ch-common.h
usr/include/colord-1/colorhug/ch-device-queue.h
usr/include/colord-1/colorhug/ch-device.h
usr/include/colord-1/colorhug/ch-hash.h
usr/include/colord-1/colorhug/ch-inhx32.h
usr/include/colord-1/colorhug/ch-math.h
usr/include/colord-1/colorhug/ch-version.h
usr/include/colord-1/colorhug/colorhug.h
usr/lib/
usr/lib/colord-plugins/
usr/lib/colord-plugins/libcd_plugin_camera.so
usr/lib/colord-plugins/libcd_plugin_sane.so
usr/lib/colord-plugins/libcd_plugin_scanner.so
usr/lib/colord-sensors/
usr/lib/colord-sensors/libcolord_sensor_argyll.so
usr/lib/colord-sensors/libcolord_sensor_colorhug.so
usr/lib/colord-sensors/libcolord_sensor_dtp94.so
usr/lib/colord-sensors/libcolord_sensor_dummy.so
usr/lib/colord-sensors/libcolord_sensor_huey.so
usr/lib/colord-sensors/libdtp94-private.so
usr/lib/colord-sensors/libhuey-private.so
usr/lib/colord-sensors/libmunki-private.so
usr/lib/colord/
usr/lib/colord/colord
usr/lib/colord/colord-sane
usr/lib/colord/colord-session
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/ColorHug-1.0.typelib
usr/lib/girepository-1.0/Colord-1.0.typelib
usr/lib/libcolord.so
usr/lib/libcolord.so.2
usr/lib/libcolord.so.2.0.3
usr/lib/libcolordprivate.so
usr/lib/libcolordprivate.so.2
usr/lib/libcolordprivate.so.2.0.3
usr/lib/libcolorhug.so
usr/lib/libcolorhug.so.2
usr/lib/libcolorhug.so.2.0.3
usr/lib/pkgconfig/
usr/lib/pkgconfig/colord.pc
usr/lib/pkgconfig/colorhug.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/colord.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/69-cd-sensors.rules
usr/lib/udev/rules.d/95-cd-devices.rules
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/colormgr
usr/share/color/
usr/share/color/icc/
usr/share/color/icc/colord/
usr/share/color/icc/colord/AdobeRGB1998.icc
usr/share/color/icc/colord/AppleRGB.icc
usr/share/color/icc/colord/BestRGB.icc
usr/share/color/icc/colord/BetaRGB.icc
usr/share/color/icc/colord/Bluish.icc
usr/share/color/icc/colord/BruceRGB.icc
usr/share/color/icc/colord/CIE-RGB.icc
usr/share/color/icc/colord/ColorMatchRGB.icc
usr/share/color/icc/colord/Crayons.icc
usr/share/color/icc/colord/DonRGB4.icc
usr/share/color/icc/colord/ECI-RGBv1.icc
usr/share/color/icc/colord/ECI-RGBv2.icc
usr/share/color/icc/colord/EktaSpacePS5.icc
usr/share/color/icc/colord/Gamma5000K.icc
usr/share/color/icc/colord/Gamma5500K.icc
usr/share/color/icc/colord/Gamma6500K.icc
usr/share/color/icc/colord/NTSC-RGB.icc
usr/share/color/icc/colord/PAL-RGB.icc
usr/share/color/icc/colord/ProPhotoRGB.icc
usr/share/color/icc/colord/SMPTE-C-RGB.icc
usr/share/color/icc/colord/SwappedRedAndGreen.icc
usr/share/color/icc/colord/WideGamutRGB.icc
usr/share/color/icc/colord/sRGB.icc
usr/share/color/icc/colord/x11-colors.icc
usr/share/colord/
usr/share/colord/cmf/
usr/share/colord/cmf/CIE1931-2deg-XYZ.cmf
usr/share/colord/cmf/CIE1964-10deg-XYZ.cmf
usr/share/colord/icons/
usr/share/colord/icons/color-munki-photo-ambient.svg
usr/share/colord/icons/color-munki-photo-attach.svg
usr/share/colord/icons/color-munki-photo-calibrate.svg
usr/share/colord/icons/color-munki-photo-projector.svg
usr/share/colord/icons/color-munki-photo-screen.svg
usr/share/colord/icons/color-munki-smile-attach.svg
usr/share/colord/icons/colorhug-attach.svg
usr/share/colord/icons/colorhug2-attach.svg
usr/share/colord/icons/dtp94-attach.svg
usr/share/colord/icons/hcfr-attach.svg
usr/share/colord/icons/huey-attach.svg
usr/share/colord/icons/i1-display3-attach.svg
usr/share/colord/icons/i1-monitor-attach.svg
usr/share/colord/icons/i1-pro-attach.svg
usr/share/colord/icons/spyder2-attach.svg
usr/share/colord/icons/spyder3-attach.svg
usr/share/colord/icons/spyder4-attach.svg
usr/share/colord/illuminant/
usr/share/colord/illuminant/CIE-A.sp
usr/share/colord/illuminant/CIE-B.sp
usr/share/colord/illuminant/CIE-C.sp
usr/share/colord/illuminant/CIE-D50.sp
usr/share/colord/illuminant/CIE-D55.sp
usr/share/colord/illuminant/CIE-D65.sp
usr/share/colord/illuminant/CIE-D93.sp
usr/share/colord/illuminant/CIE-E.sp
usr/share/colord/illuminant/CIE-F1.sp
usr/share/colord/illuminant/CIE-F10.sp
usr/share/colord/illuminant/CIE-F11.sp
usr/share/colord/illuminant/CIE-F12.sp
usr/share/colord/illuminant/CIE-F2.sp
usr/share/colord/illuminant/CIE-F3.sp
usr/share/colord/illuminant/CIE-F4.sp
usr/share/colord/illuminant/CIE-F5.sp
usr/share/colord/illuminant/CIE-F6.sp
usr/share/colord/illuminant/CIE-F7.sp
usr/share/colord/illuminant/CIE-F8.sp
usr/share/colord/illuminant/CIE-F9.sp
usr/share/colord/ref/
usr/share/colord/ref/CIE-1986-daylight-SPD.cmf
usr/share/colord/ref/CIE-TCS.sp
usr/share/colord/ti1/
usr/share/colord/ti1/display-long.ti1
usr/share/colord/ti1/display-normal.ti1
usr/share/colord/ti1/display-short.ti1
usr/share/colord/ti1/printer-long.ti1
usr/share/colord/ti1/printer-normal.ti1
usr/share/colord/ti1/printer-short.ti1
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.ColorHelper.xml
usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Device.xml
usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Profile.xml
usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Sensor.xml
usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.xml
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.freedesktop.ColorHelper.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.ColorManager.service
usr/share/gir-1.0/
usr/share/gir-1.0/ColorHug-1.0.gir
usr/share/gir-1.0/Colord-1.0.gir
usr/share/glib-2.0/
usr/share/glib-2.0/schemas/
usr/share/glib-2.0/schemas/org.freedesktop.ColorHelper.gschema.xml
usr/share/locale/
usr/share/locale/as/
usr/share/locale/as/LC_MESSAGES/
usr/share/locale/as/LC_MESSAGES/colord.mo
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/colord.mo
usr/share/locale/bg_BG/
usr/share/locale/bg_BG/LC_MESSAGES/
usr/share/locale/bg_BG/LC_MESSAGES/colord.mo
usr/share/locale/bn_IN/
usr/share/locale/bn_IN/LC_MESSAGES/
usr/share/locale/bn_IN/LC_MESSAGES/colord.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/colord.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/colord.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/colord.mo
usr/share/locale/da_DK/
usr/share/locale/da_DK/LC_MESSAGES/
usr/share/locale/da_DK/LC_MESSAGES/colord.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/colord.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/colord.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/colord.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/colord.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/colord.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/colord.mo
usr/share/locale/fa/
usr/share/locale/fa/LC_MESSAGES/
usr/share/locale/fa/LC_MESSAGES/colord.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/colord.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/colord.mo
usr/share/locale/ga/
usr/share/locale/ga/LC_MESSAGES/
usr/share/locale/ga/LC_MESSAGES/colord.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/colord.mo
usr/share/locale/gu/
usr/share/locale/gu/LC_MESSAGES/
usr/share/locale/gu/LC_MESSAGES/colord.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/colord.mo
usr/share/locale/hi/
usr/share/locale/hi/LC_MESSAGES/
usr/share/locale/hi/LC_MESSAGES/colord.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/colord.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/colord.mo
usr/share/locale/ia/
usr/share/locale/ia/LC_MESSAGES/
usr/share/locale/ia/LC_MESSAGES/colord.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/colord.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/colord.mo
usr/share/locale/it_IT/
usr/share/locale/it_IT/LC_MESSAGES/
usr/share/locale/it_IT/LC_MESSAGES/colord.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/colord.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/colord.mo
usr/share/locale/kk/
usr/share/locale/kk/LC_MESSAGES/
usr/share/locale/kk/LC_MESSAGES/colord.mo
usr/share/locale/kn/
usr/share/locale/kn/LC_MESSAGES/
usr/share/locale/kn/LC_MESSAGES/colord.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/colord.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/colord.mo
usr/share/locale/lv/
usr/share/locale/lv/LC_MESSAGES/
usr/share/locale/lv/LC_MESSAGES/colord.mo
usr/share/locale/ml/
usr/share/locale/ml/LC_MESSAGES/
usr/share/locale/ml/LC_MESSAGES/colord.mo
usr/share/locale/mr/
usr/share/locale/mr/LC_MESSAGES/
usr/share/locale/mr/LC_MESSAGES/colord.mo
usr/share/locale/ms/
usr/share/locale/ms/LC_MESSAGES/
usr/share/locale/ms/LC_MESSAGES/colord.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/colord.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/colord.mo
usr/share/locale/nn/
usr/share/locale/nn/LC_MESSAGES/
usr/share/locale/nn/LC_MESSAGES/colord.mo
usr/share/locale/or/
usr/share/locale/or/LC_MESSAGES/
usr/share/locale/or/LC_MESSAGES/colord.mo
usr/share/locale/pa/
usr/share/locale/pa/LC_MESSAGES/
usr/share/locale/pa/LC_MESSAGES/colord.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/colord.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/colord.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/colord.mo
usr/share/locale/pt_PT/
usr/share/locale/pt_PT/LC_MESSAGES/
usr/share/locale/pt_PT/LC_MESSAGES/colord.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/colord.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/colord.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/colord.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/colord.mo
usr/share/locale/sq/
usr/share/locale/sq/LC_MESSAGES/
usr/share/locale/sq/LC_MESSAGES/colord.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/colord.mo
usr/share/locale/sr@latin/
usr/share/locale/sr@latin/LC_MESSAGES/
usr/share/locale/sr@latin/LC_MESSAGES/colord.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/colord.mo
usr/share/locale/ta/
usr/share/locale/ta/LC_MESSAGES/
usr/share/locale/ta/LC_MESSAGES/colord.mo
usr/share/locale/te/
usr/share/locale/te/LC_MESSAGES/
usr/share/locale/te/LC_MESSAGES/colord.mo
usr/share/locale/th/
usr/share/locale/th/LC_MESSAGES/
usr/share/locale/th/LC_MESSAGES/colord.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/colord.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/colord.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/colord.mo
usr/share/locale/zh_HK/
usr/share/locale/zh_HK/LC_MESSAGES/
usr/share/locale/zh_HK/LC_MESSAGES/colord.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/colord.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cd-create-profile.1.gz
usr/share/man/man1/cd-fix-profile.1.gz
usr/share/man/man1/cd-it8.1.gz
usr/share/man/man1/colormgr.1.gz
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.freedesktop.color.policy
usr/share/vala/
usr/share/vala/vapi/
usr/share/vala/vapi/colord.vapi
var/
var/lib/
var/lib/colord/
var/lib/colord/icc/

