# var/lib/pacman/local/geoclue2-2.1.10-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.freedesktop.GeoClue2.Agent.conf
etc/dbus-1/system.d/org.freedesktop.GeoClue2.conf
etc/geoclue/
etc/geoclue/geoclue.conf
usr/
usr/lib/
usr/lib/geoclue2/
usr/lib/geoclue2/geoclue
usr/lib/geoclue2/geoclue-2.0/
usr/lib/geoclue2/geoclue-2.0/demos/
usr/lib/geoclue2/geoclue-2.0/demos/where-am-i
usr/lib/pkgconfig/
usr/lib/pkgconfig/geoclue-2.0.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/geoclue.service
usr/share/
usr/share/applications/
usr/share/applications/geoclue-demo-agent.desktop
usr/share/applications/geoclue-where-am-i.desktop
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Agent.xml
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.GeoClue2.service

