# var/lib/pacman/local/avahi-0.6.31-15/files
%FILES%
etc/
etc/avahi/
etc/avahi/avahi-autoipd.action
etc/avahi/avahi-daemon.conf
etc/avahi/avahi-dnsconfd.action
etc/avahi/hosts
etc/avahi/services/
etc/avahi/services/sftp-ssh.service
etc/avahi/services/ssh.service
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/avahi-dbus.conf
usr/
usr/bin/
usr/bin/avahi-autoipd
usr/bin/avahi-bookmarks
usr/bin/avahi-browse
usr/bin/avahi-browse-domains
usr/bin/avahi-daemon
usr/bin/avahi-discover
usr/bin/avahi-discover-standalone
usr/bin/avahi-dnsconfd
usr/bin/avahi-publish
usr/bin/avahi-publish-address
usr/bin/avahi-publish-service
usr/bin/avahi-resolve
usr/bin/avahi-resolve-address
usr/bin/avahi-resolve-host-name
usr/bin/avahi-set-host-name
usr/bin/bshell
usr/bin/bssh
usr/bin/bvnc
usr/include/
usr/include/avahi-client/
usr/include/avahi-client/client.h
usr/include/avahi-client/lookup.h
usr/include/avahi-client/publish.h
usr/include/avahi-common/
usr/include/avahi-common/address.h
usr/include/avahi-common/alternative.h
usr/include/avahi-common/cdecl.h
usr/include/avahi-common/defs.h
usr/include/avahi-common/domain.h
usr/include/avahi-common/error.h
usr/include/avahi-common/gccmacro.h
usr/include/avahi-common/llist.h
usr/include/avahi-common/malloc.h
usr/include/avahi-common/rlist.h
usr/include/avahi-common/simple-watch.h
usr/include/avahi-common/strlst.h
usr/include/avahi-common/thread-watch.h
usr/include/avahi-common/timeval.h
usr/include/avahi-common/watch.h
usr/include/avahi-compat-howl/
usr/include/avahi-compat-howl/corby/
usr/include/avahi-compat-howl/corby/buffer.h
usr/include/avahi-compat-howl/corby/channel.h
usr/include/avahi-compat-howl/corby/corby.h
usr/include/avahi-compat-howl/corby/message.h
usr/include/avahi-compat-howl/corby/object.h
usr/include/avahi-compat-howl/corby/orb.h
usr/include/avahi-compat-howl/discovery/
usr/include/avahi-compat-howl/discovery/discovery.h
usr/include/avahi-compat-howl/discovery/text_record.h
usr/include/avahi-compat-howl/howl.h
usr/include/avahi-compat-howl/rendezvous/
usr/include/avahi-compat-howl/rendezvous/rendezvous.h
usr/include/avahi-compat-howl/rendezvous/text_record.h
usr/include/avahi-compat-howl/salt/
usr/include/avahi-compat-howl/salt/address.h
usr/include/avahi-compat-howl/salt/debug.h
usr/include/avahi-compat-howl/salt/interface.h
usr/include/avahi-compat-howl/salt/platform.h
usr/include/avahi-compat-howl/salt/salt.h
usr/include/avahi-compat-howl/salt/signal.h
usr/include/avahi-compat-howl/salt/socket.h
usr/include/avahi-compat-howl/salt/time.h
usr/include/avahi-compat-libdns_sd/
usr/include/avahi-compat-libdns_sd/dns_sd.h
usr/include/avahi-core/
usr/include/avahi-core/core.h
usr/include/avahi-core/log.h
usr/include/avahi-core/lookup.h
usr/include/avahi-core/publish.h
usr/include/avahi-core/rr.h
usr/include/avahi-glib/
usr/include/avahi-glib/glib-malloc.h
usr/include/avahi-glib/glib-watch.h
usr/include/avahi-gobject/
usr/include/avahi-gobject/ga-client.h
usr/include/avahi-gobject/ga-entry-group.h
usr/include/avahi-gobject/ga-enums.h
usr/include/avahi-gobject/ga-error.h
usr/include/avahi-gobject/ga-record-browser.h
usr/include/avahi-gobject/ga-service-browser.h
usr/include/avahi-gobject/ga-service-resolver.h
usr/include/avahi-qt4/
usr/include/avahi-qt4/qt-watch.h
usr/include/avahi-ui/
usr/include/avahi-ui/avahi-ui.h
usr/include/dns_sd.h
usr/include/howl
usr/lib/
usr/lib/avahi/
usr/lib/avahi/service-types.db
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/Avahi-0.6.typelib
usr/lib/girepository-1.0/AvahiCore-0.6.typelib
usr/lib/libavahi-client.so
usr/lib/libavahi-client.so.3
usr/lib/libavahi-client.so.3.2.9
usr/lib/libavahi-common.so
usr/lib/libavahi-common.so.3
usr/lib/libavahi-common.so.3.5.3
usr/lib/libavahi-core.so
usr/lib/libavahi-core.so.7
usr/lib/libavahi-core.so.7.0.2
usr/lib/libavahi-glib.so
usr/lib/libavahi-glib.so.1
usr/lib/libavahi-glib.so.1.0.2
usr/lib/libavahi-gobject.so
usr/lib/libavahi-gobject.so.0
usr/lib/libavahi-gobject.so.0.0.4
usr/lib/libavahi-qt4.so
usr/lib/libavahi-qt4.so.1
usr/lib/libavahi-qt4.so.1.0.2
usr/lib/libavahi-ui-gtk3.so
usr/lib/libavahi-ui-gtk3.so.0
usr/lib/libavahi-ui-gtk3.so.0.1.4
usr/lib/libavahi-ui.so
usr/lib/libavahi-ui.so.0
usr/lib/libavahi-ui.so.0.1.4
usr/lib/libdns_sd.so
usr/lib/libdns_sd.so.1
usr/lib/libdns_sd.so.1.0.0
usr/lib/libhowl.so
usr/lib/libhowl.so.0
usr/lib/libhowl.so.0.0.0
usr/lib/mono/
usr/lib/mono/avahi-sharp/
usr/lib/mono/avahi-sharp/avahi-sharp.dll
usr/lib/mono/avahi-ui-sharp/
usr/lib/mono/avahi-ui-sharp/avahi-ui-sharp.dll
usr/lib/mono/gac/
usr/lib/mono/gac/avahi-sharp/
usr/lib/mono/gac/avahi-sharp/1.0.0.0__4d116c78973743f5/
usr/lib/mono/gac/avahi-sharp/1.0.0.0__4d116c78973743f5/avahi-sharp.dll
usr/lib/mono/gac/avahi-sharp/1.0.0.0__4d116c78973743f5/avahi-sharp.dll.config
usr/lib/mono/gac/avahi-sharp/1.0.0.0__4d116c78973743f5/avahi-sharp.dll.mdb
usr/lib/mono/gac/avahi-ui-sharp/
usr/lib/mono/gac/avahi-ui-sharp/0.0.0.0__4d116c78973743f5/
usr/lib/mono/gac/avahi-ui-sharp/0.0.0.0__4d116c78973743f5/avahi-ui-sharp.dll
usr/lib/mono/gac/avahi-ui-sharp/0.0.0.0__4d116c78973743f5/avahi-ui-sharp.dll.config
usr/lib/mono/gac/avahi-ui-sharp/0.0.0.0__4d116c78973743f5/avahi-ui-sharp.dll.mdb
usr/lib/pkgconfig/
usr/lib/pkgconfig/avahi-client.pc
usr/lib/pkgconfig/avahi-compat-howl.pc
usr/lib/pkgconfig/avahi-compat-libdns_sd.pc
usr/lib/pkgconfig/avahi-core.pc
usr/lib/pkgconfig/avahi-glib.pc
usr/lib/pkgconfig/avahi-gobject.pc
usr/lib/pkgconfig/avahi-qt4.pc
usr/lib/pkgconfig/avahi-sharp.pc
usr/lib/pkgconfig/avahi-ui-gtk3.pc
usr/lib/pkgconfig/avahi-ui-sharp.pc
usr/lib/pkgconfig/avahi-ui.pc
usr/lib/pkgconfig/howl.pc
usr/lib/python2.7/
usr/lib/python2.7/site-packages/
usr/lib/python2.7/site-packages/avahi/
usr/lib/python2.7/site-packages/avahi/ServiceTypeDatabase.py
usr/lib/python2.7/site-packages/avahi/__init__.py
usr/lib/python2.7/site-packages/avahi/__init__.pyc
usr/lib/python2.7/site-packages/avahi/__init__.pyo
usr/lib/python2.7/site-packages/avahi_discover/
usr/lib/python2.7/site-packages/avahi_discover/__init__.py
usr/lib/python2.7/site-packages/avahi_discover/__init__.pyc
usr/lib/python2.7/site-packages/avahi_discover/__init__.pyo
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/avahi-daemon.service
usr/lib/systemd/system/avahi-daemon.socket
usr/lib/systemd/system/avahi-dnsconfd.service
usr/lib/systemd/system/dbus-org.freedesktop.Avahi.service
usr/share/
usr/share/applications/
usr/share/applications/avahi-discover.desktop
usr/share/applications/bssh.desktop
usr/share/applications/bvnc.desktop
usr/share/avahi/
usr/share/avahi/avahi-service.dtd
usr/share/avahi/interfaces/
usr/share/avahi/interfaces/avahi-discover.ui
usr/share/avahi/service-types
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.AddressResolver.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.DomainBrowser.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.EntryGroup.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.HostNameResolver.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.RecordBrowser.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.Server.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceBrowser.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceResolver.xml
usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceTypeBrowser.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.Avahi.service
usr/share/gir-1.0/
usr/share/gir-1.0/Avahi-0.6.gir
usr/share/gir-1.0/AvahiCore-0.6.gir
usr/share/locale/
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/avahi.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/avahi.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/avahi.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/avahi.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/avahi.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/avahi.mo
usr/share/locale/en_AU/
usr/share/locale/en_AU/LC_MESSAGES/
usr/share/locale/en_AU/LC_MESSAGES/avahi.mo
usr/share/locale/en_CA/
usr/share/locale/en_CA/LC_MESSAGES/
usr/share/locale/en_CA/LC_MESSAGES/avahi.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/avahi.mo
usr/share/locale/en_NZ/
usr/share/locale/en_NZ/LC_MESSAGES/
usr/share/locale/en_NZ/LC_MESSAGES/avahi.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/avahi.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/avahi.mo
usr/share/locale/fo/
usr/share/locale/fo/LC_MESSAGES/
usr/share/locale/fo/LC_MESSAGES/avahi.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/avahi.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/avahi.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/avahi.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/avahi.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/avahi.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/avahi.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/avahi.mo
usr/share/locale/ms/
usr/share/locale/ms/LC_MESSAGES/
usr/share/locale/ms/LC_MESSAGES/avahi.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/avahi.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/avahi.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/avahi.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/avahi.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/avahi.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/avahi.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/avahi.mo
usr/share/locale/sr@latin/
usr/share/locale/sr@latin/LC_MESSAGES/
usr/share/locale/sr@latin/LC_MESSAGES/avahi.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/avahi.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/avahi.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/avahi.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/avahi.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/avahi-bookmarks.1.gz
usr/share/man/man1/avahi-browse-domains.1.gz
usr/share/man/man1/avahi-browse.1.gz
usr/share/man/man1/avahi-discover.1.gz
usr/share/man/man1/avahi-publish-address.1.gz
usr/share/man/man1/avahi-publish-service.1.gz
usr/share/man/man1/avahi-publish.1.gz
usr/share/man/man1/avahi-resolve-address.1.gz
usr/share/man/man1/avahi-resolve-host-name.1.gz
usr/share/man/man1/avahi-resolve.1.gz
usr/share/man/man1/avahi-set-host-name.1.gz
usr/share/man/man1/bssh.1.gz
usr/share/man/man1/bvnc.1.gz
usr/share/man/man5/
usr/share/man/man5/avahi-daemon.conf.5.gz
usr/share/man/man5/avahi.hosts.5.gz
usr/share/man/man5/avahi.service.5.gz
usr/share/man/man8/
usr/share/man/man8/avahi-autoipd.8.gz
usr/share/man/man8/avahi-autoipd.action.8.gz
usr/share/man/man8/avahi-daemon.8.gz
usr/share/man/man8/avahi-dnsconfd.8.gz
usr/share/man/man8/avahi-dnsconfd.action.8.gz

%BACKUP%
etc/avahi/hosts	186990ae1edac95a88dbef6a36a07716
etc/avahi/avahi-daemon.conf	e8f43bcc5d819ab0843cc55f21eaaa9d
etc/avahi/services/ssh.service	259360774d926043ba6674514adac3d2
etc/avahi/services/sftp-ssh.service	c72d539947ec6966acbc3b87b87b5b83
usr/lib/avahi/service-types.db	b345fc41655eacbdad7d495624082ff5
usr/share/avahi/service-types	f228c5a9125cf37536c86d5d293997d7

