# var/lib/pacman/local/geoclue2-2.4.1-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.freedesktop.GeoClue2.Agent.conf
etc/dbus-1/system.d/org.freedesktop.GeoClue2.conf
etc/geoclue/
etc/geoclue/geoclue.conf
usr/
usr/include/
usr/include/libgeoclue-2.0/
usr/include/libgeoclue-2.0/gclue-client.h
usr/include/libgeoclue-2.0/gclue-enum-types.h
usr/include/libgeoclue-2.0/gclue-enums.h
usr/include/libgeoclue-2.0/gclue-helpers.h
usr/include/libgeoclue-2.0/gclue-location.h
usr/include/libgeoclue-2.0/gclue-manager.h
usr/include/libgeoclue-2.0/gclue-simple.h
usr/include/libgeoclue-2.0/geoclue.h
usr/lib/
usr/lib/geoclue2/
usr/lib/geoclue2/geoclue
usr/lib/geoclue2/geoclue-2.0/
usr/lib/geoclue2/geoclue-2.0/demos/
usr/lib/geoclue2/geoclue-2.0/demos/where-am-i
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/Geoclue-2.0.typelib
usr/lib/libgeoclue-2.so
usr/lib/libgeoclue-2.so.0
usr/lib/libgeoclue-2.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/geoclue-2.0.pc
usr/lib/pkgconfig/libgeoclue-2.0.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/geoclue.service
usr/share/
usr/share/applications/
usr/share/applications/geoclue-demo-agent.desktop
usr/share/applications/geoclue-where-am-i.desktop
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Agent.xml
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Client.xml
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Location.xml
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Manager.xml
usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.GeoClue2.service
usr/share/gir-1.0/
usr/share/gir-1.0/Geoclue-2.0.gir

%BACKUP%
etc/geoclue/geoclue.conf	52d16f4780e084aa7041e3f2a6d16df6

