# var/lib/pacman/local/systemd-readahead-216-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/systemd-readahead-collect.service
usr/lib/systemd/system/systemd-readahead-done.service
usr/lib/systemd/system/systemd-readahead-done.timer
usr/lib/systemd/system/systemd-readahead-drop.service
usr/lib/systemd/system/systemd-readahead-replay.service
usr/lib/systemd/systemd-readahead
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/systemd-readahead-collect.service.8.gz
usr/share/man/man8/systemd-readahead-done.service.8.gz
usr/share/man/man8/systemd-readahead-done.timer.8.gz
usr/share/man/man8/systemd-readahead-replay.service.8.gz
usr/share/man/man8/systemd-readahead.8.gz

